# package/dnspod-OpenWRT/files/etc/hotplug.d/iface/99-dnspod-updater
#!/bin/sh
[ "$ACTION" = ifup ] && {
    logger -p6 -t DNSPOD "dnspod-updater is running."
    /usr/sbin/dnspod 
}
